# conf/test/supervisor/programs/xxb_send.conf
[program:xxb_send]
numprocs = 1
process_name = bb-%(process_num)s
command = /usr/local/php/bin/php /mnt/htdocs/cli.php request_uri=/queue/xxbSend env=test &
autorestart = true
redirect_stderr = true
stdout_logfile = /mnt/data/log/xxb_send.log
stderr_logfile = /mnt/data/log/xxb_send-error.log
